# services/gmail/gmail-backup@.timer
[Unit]
Description=Run Gmail mailbox backup monthly on the 1st at 15:00

[Timer]
OnCalendar=*-*-01 15:00:00
Persistent=true

[Install]
WantedBy=timers.target
