# services/gmail/gmail-backup@.service
[Unit]
Description=Gmail mailbox backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%i
LoadCredential=gmail-user:/etc/credstore/gmail-user
LoadCredential=gmail-app-password:/etc/credstore/gmail-app-password
LoadCredential=encryption-pub-key:/etc/credstore/gmail-encryption-pub-key
ExecStart=/usr/local/libexec/dotfiles/gmail-backup.sh

# Sandboxing (common directives appended from lib/sandbox.conf at install time)
BindPaths=/home/%i/Backups/Email
ReadWritePaths=/home/%i/Backups/Email
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
