# services/storage/s3-backup@.timer
[Unit]
Description=Run S3 backup monthly on the 1st at 23:00

[Timer]
OnCalendar=*-*-01 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
